# systemd/git-remote-push.timer
[Unit]
Description=Git Remote Push Timer

[Timer]
OnUnitActiveSec=4h
RandomizedDelaySec=120m
Persistent=true

[Install]
WantedBy=timers.target
